# scripts/systemd/vps-metrics-exporter.timer
[Unit]
Description=Run VPS Metrics Exporter every 60 seconds

[Timer]
OnCalendar=*:0/1
Persistent=true

[Install]
WantedBy=timers.target
